# live/collector/muretto-live.service
# Unit systemd del collettore live — installazione:
#   sudo cp live/collector/muretto-live.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now muretto-live
# Log:  journalctl -u muretto-live -f
[Unit]
Description=Muretto — collettore live timing F1 (Fase 2)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=muretto
WorkingDirectory=/home/muretto/muretto
# PITLANE_ARGS (es. "--pitlane data/live_derived/pitlane_spa.json") viene
# dal file env sotto: si cambia per weekend/circuito senza toccare l'unit
# (vedi live/RUNBOOK_WEEKEND.md). File assente = niente in_pit geometrico.
EnvironmentFile=-/home/muretto/.muretto-live.env
ExecStart=/home/muretto/muretto/.venv-live/bin/python live/collector/collector.py --out-dir /home/muretto/muretto/data/live_raw $PITLANE_ARGS
# Il daemon non esce mai per progetto; il restart e' solo cintura di
# sicurezza (crash imprevisti). Mai piu' aggressivo di cosi'.
Restart=always
# 30s (era 5s): la guardia sui CONNECT vive NEL processo, quindi un
# crash-loop la aggirerebbe — a 5s farebbe 12 CONNECT/minuto, sopra la
# soglia OpenF1 di 10 disconnessioni/minuto (blocco 10 min; mail OpenF1
# del 21/07/2026, vedi REPORT_MQTT_INGRESS.md). A 30s: 2/minuto.
RestartSec=30

[Install]
WantedBy=multi-user.target
